# new-voice.service
[Unit]
Description=get obj from PIPE,and translate to voice by baidu api.

[Service]
WorkingDirectory=/NewsFMAPI/voice/
ExecStart=./run.sh